# etc/uptimed-bootid.xml.in
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type="manifest" name="uptimed-bootid">

<service
	name="site/uptimed-bootid"
	type="service"
	version="1">

	<create_default_instance enabled="false" />

	<single_instance />

	<dependency name="fs-local"
		grouping="require_all"
		restart_on="none"
		type="service">
		<service_fmri value="svc:/system/filesystem/local" />
	</dependency>

        <dependency
                name="milestone"
                grouping="require_all"
                restart_on="none"
                type="service">
                <service_fmri value="svc:/milestone/single-user" />
        </dependency>

	<exec_method
		type="method"
		name="start"
		exec="@prefix@/sbin/uptimed -b"
		timeout_seconds="-1" />

	<exec_method
		type="method"
		name="stop"
		exec=":true"
		timeout_seconds="-1" />

	<property_group name="startd" type="framework">
		<propval name="duration" type="astring" value="transient" />
	</property_group>

	<stability value="Unstable" />

	<template>
		<common_name>
			<loctext xml:lang="C">
				Create boot ID for uptimed(8) on system boots up
			</loctext>
		</common_name>
		<documentation>
			<manpage title="uptimed" section="8" manpath="@prefix@/share/man" />
		</documentation>
	</template>

</service>

</service_bundle>
