# etc/uptimed.service.in
[Unit]
Description=Uptime record tracking daemon
Documentation=man:uptimed(8) man:uprecords(1)

[Service]
Type=notify
ExecStart=@prefix@/sbin/uptimed -f
Restart=on-failure

[Install]
WantedBy=multi-user.target
